# systemd/user/rob-config-drift-minimal.service
[Unit]
Description=Rob Config Drift Check
After=default.target

[Service]
Type=oneshot
ExecStart=/home/yama/ws/scripts/config-drift-check.sh
WorkingDirectory=/home/yama/ws
Nice=5
